# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:wechat-bot]
command=npx wechatbot-webhook
environment=LOGIN_API_TOKEN="trip-bot-2026",RECVD_MSG_API="http://localhost:7860/api/wechat/callback"
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=10
startsecs=5
priority=1

[program:flask]
command=gunicorn -b 0.0.0.0:7860 --workers 1 --threads 2 --timeout 120 app:app
directory=/app
environment=DATABASE="/app/data/trip.db",BOT_BASE_URL="http://localhost:3001",BOT_TOKEN="trip-bot-2026",BOT_LOGIN_PREFIX="/bot"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=5
startsecs=3
priority=10
